# 各种配置文件/supervisor/deploy.ini
[program:deploy] 
command=/usr/bin/python3 /root/worker/deploy.py
director=/root/worker
autorestart=true          
redirect_stderr=true      
stdout_logfile = /var/log/supervisord/deploy.log
loglevel=info
numprocs=5                   
process_name=%(program_name)s_%(process_num)02d; 
user=root
