# 各种配置文件/supervisor/fabric.ini
[program:fabric]
process_name=%(program_name)s_%(process_num)01d; 
command=/usr/bin/python3 /root/pyweb/application.py --port=%(process_num)d --debug
directory=/root/pyweb
autostart=true
autorestart=true
startsecs=10
startretries=20
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisord/fabric.log
stderr_logfile=/var/log/supervisord/fabric.log
numprocs=8
numprocs_start=8880

